# experimental/production_poc/deploy/production-poc-monitor.timer
[Unit]
Description=Run the production emergency recovery PoC every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=production-poc-monitor.service

[Install]
WantedBy=timers.target
